# files/boot.unit.service
[Unit]
Description=Verific configuratia unPi pentru beneficiari

[Service]
Type=oneshot
ExecStart=/bin/dash /root/.unpi/boot.sh
