# files/boot.unit.timer
[Unit]
Description=Verific configuratia unPi (timer)

[Timer]
OnUnitActiveSec=30min
OnBootSec=60s

[Install]
WantedBy=timers.target
